# init/roslaunch.service
[Unit]
Description=roslaunch

[Service]
ExecStart=/home/pluto/bin/doroslaunch.sh

[Install]
WantedBy=default.target
